# bluealsa-handsfree-control.socket.in
[Socket]
ListenFIFO=%t/bluealsa-handsfree/control
SocketMode=0622
RemoveOnStop=true

[Install]
WantedBy=sockets.target
